# etc/crontab
5 * * * * root /usr/sbin/logrotate /etc/terragen/logrotate.conf
