# etc/crontab
5 * * * * root logrotate /etc/terragen/logrotate.conf
